# api.conf
[program:api]
command=/srv/cloudplayer/api/bin/api --config=aws.py
user=cloudplayer
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
stopsignal=TERM
stopwaitsecs=10
