# supervisord.conf
[program:api]
command=/var/www/api/bin/api --config=aws.py
user=www-data
numprocs=1
autostart=true
autorestart=true
startsecs=10
startretries=3
stopsignal=TERM
stopwaitsecs=10
